# deploy/ytdlp-pot.service
[Unit]
Description=yt-dlp PO Token provider (bgutil) for YouTube bypass
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@BOT_USER@
WorkingDirectory=@BOT_HOME@/bgutil-ytdlp-pot-provider/server
ExecStart=/usr/bin/node build/main.js --port @POT_PORT@
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ytdlp-pot

[Install]
WantedBy=multi-user.target
